# scripts/com.reyansh.erp.daily-deploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.reyansh.erp.daily-deploy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/abhishekjain/Desktop/reyansh-erp-new/scripts/daily-deploy.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/abhishekjain/Desktop/reyansh-erp-new/scripts/daily-deploy.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/abhishekjain/Desktop/reyansh-erp-new/scripts/daily-deploy.err.log</string>
</dict>
</plist>
